# cacti.crontab
MAILTO=root

*/5 * * * * cacti umask 022; exec /usr/sbin/cacti-poller > /dev/null

# vim:ft=crontab
