# cacti.crontab
MAILTO=root

*/5 * * * * http umask 022; /usr/sbin/cacti-poller

# vim:ft=crontab
